# services/launchd/com.jinchiwei.jkw-obs-mcp.daily-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jinchiwei.jkw-obs-mcp.daily-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/python</string>
        <string>-m</string>
        <string>jkw_obs_mcp.triggers.daily_review_runner</string>
    </array>

    <!-- Fires every 5 min while awake. Paused during sleep; resumes on wake.
         RunAtLoad=true also fires once at session login (covers full restarts). -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Don't start a second instance if a previous run is still going. -->
    <key>AbandonProcessGroup</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jinchiwei/Library/Logs/com.jinchiwei.jkw-obs-mcp.daily-review.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jinchiwei/Library/Logs/com.jinchiwei.jkw-obs-mcp.daily-review.err</string>
</dict>
</plist>

